# front.conf
[program:yummy]
command=npm start
directory=/home/ubuntu/yummy-react
user=ubuntu
autostart=true
autorestart=unexpected
stdout_logfile=/home/ubuntu/react.log
stderr_logfile=/home/ubuntu/react.err.log
